# app/client/files/etc/systemd/system/displayer-conf-watchdog.service
[Unit]
Description=displayer-conf-watchdog

[Service]
ExecStart=/opt/displayer/displayer-conf-watchdog

[Install]
WantedBy=multi-user.target
